# deploy/elephant.service.d/neo4j.conf
# Drop-in: make elephant wait for the elephant-neo4j unit on boot.
# Installed to /etc/systemd/system/elephant.service.d/neo4j.conf
[Unit]
After=elephant-neo4j.service
Requires=elephant-neo4j.service
